# lab_17/files/units/borg-run.timer.j2
[Unit]
Description=Borg Backup
Requires=borg-run.service

[Timer]
Unit=borg-run.service
OnActiveSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
